# sermon-broadcaster.service
[Unit]
Description=Sermon Broadcaster
After=network.target

[Service]
Type=simple
User=lherzog
WorkingDirectory=/opt/sermon-broadcaster
EnvironmentFile=/opt/sermon-broadcaster/.env
# Transcription backend env vars (operative values live in .env — these are
# documentation/overrides; uncommenting here would override .env, so leave
# the toggle in .env unless you intend to pin it at the unit level):
#   TRANSCRIBE_BACKEND=openai        # openai (default) | local | faster-whisper
#   WHISPER_LOCAL_URL=http://10.0.0.50:5005   # mini static IP (not .local)
#   WHISPER_LOCAL_TIMEOUT=600        # seconds; ~75-min sermon = 5-8 min call
# Keep one worker because in-process analysis/render threads and import-time job
# reconciliation are not safe across multiple worker processes.
ExecStart=/opt/sermon-broadcaster/venv/bin/gunicorn --workers 1 --threads 4 --timeout 0 --bind 0.0.0.0:5003 app:app
Restart=on-failure
RestartSec=5
# Large file uploads need time
TimeoutStartSec=30
TimeoutStopSec=30
# Prevent OOM from Whisper if using local transcription
MemoryMax=8G

[Install]
WantedBy=multi-user.target
